# salt/salt/ipython/files/notebook.conf
{% set environment = [] %}
{% do environment.append('USER="dsb"') %}
{% do environment.append('IPYTHONDIR="/home/dsb/.ipython"') %}
{% do environment.append('JAVA_HOME="/usr/java/default"') %}
{% do environment.append('SPARK_HOME="/usr/lib/spark"') %}
{% do environment.append('PYTHONPATH="/usr/lib/spark/python:/usr/lib/spark/python/lib/py4j-0.8.1-src.zip:$PYTHONPATH"') %}
{% do environment.append('MESOS_NATIVE_LIBRARY="/usr/local/lib/libmesos.so"') %}
{% do environment.append('SPARK_EXECUTOR_URI="hdfs://ip-10-187-18-185:8020/tmp/spark-1.0.2-bin-hadoop2.tgz"') %}
{% do environment.append('MASTER="zk://ip-10-187-18-185:2181/mesos"') %}
{% do environment.append('PYSPARK_SUBMIT_ARGS="--master mesos://daniel-master:5050"') %}

[program:notebook]
command=/home/dsb/envs/base/bin/ipython notebook --ip=0.0.0.0 --port=8888 --no-browser
user=dsb
numprocs=1
directory=/home/dsb/notebooks
autostart=false
autorestart=true
stdout_logfile=/var/log/ipython/notebook.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
environment={{ environment | join(',') }}
